# config/launchd/ai.earnedout.overnight-search.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.earnedout.overnight-search</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/biffreybraxton/published-listing-search/run-overnight-search.sh</string>
    </array>

    <!-- Nightly, early morning: 02:37 local time (off-minute on purpose). -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>2</integer>
        <key>Minute</key>
        <integer>37</integer>
    </dict>

    <!-- Do not run when the agent is (re)loaded — only on the calendar trigger. -->
    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/biffreybraxton/published-listing-search/output/logs/launchd.overnight-search.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/biffreybraxton/published-listing-search/output/logs/launchd.overnight-search.err.log</string>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
